# system_files/usr/lib/systemd/system/niriblue-flatpak-setup.service
[Unit]
Description=Install Flathub apps on first boot
After=network-online.target
Wants=network-online.target
ConditionPathExists=!/var/lib/niriblue/flatpak-setup.done

[Service]
Type=oneshot
ExecStart=/usr/libexec/niriblue-flatpak-setup
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
