# system_files/usr/lib/systemd/system/niriblue-firstboot.service
[Unit]
Description=niriblue first-boot setup (gated, with on-screen progress)
# Run while the Plymouth splash is still up and BEFORE the greeter, so the desktop does
# not appear until first-boot installs are done. network-online so the install steps can
# reach the network; the splash stays because we order before plymouth-quit*.
After=network-online.target plymouth-start.service
Wants=network-online.target
Before=plymouth-quit.service plymouth-quit-wait.service greetd.service display-manager.service plasmalogin.service
ConditionPathExists=!/var/lib/niriblue/firstboot.done

[Service]
Type=oneshot
ExecStart=/usr/libexec/niriblue-firstboot
RemainAfterExit=yes
# Hard backstop: the orchestrator self-caps each step, but never let a stuck boot trap
# the user behind the splash forever.
TimeoutStartSec=1800
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
