# system_files/usr/lib/systemd/system/niriblue-sysext-setup.service
[Unit]
Description=Bootstrap sysexts-manager + niriblue sysexts (vscode, tailscale) on first boot
# Ordered after the first-boot orchestrator (see niriblue-nix-setup.service for why).
After=network-online.target systemd-sysext.service niriblue-firstboot.service
Wants=network-online.target
ConditionPathExists=!/var/lib/niriblue/sysext-setup.done

[Service]
Type=oneshot
ExecStart=/usr/libexec/niriblue-sysext-setup
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
