# system_files/usr/lib/systemd/system/niriblue-flatpak-setup.service
[Unit]
Description=Install Flathub apps on first boot
# Ordered after the first-boot orchestrator (see niriblue-nix-setup.service for why).
After=network-online.target niriblue-firstboot.service
Wants=network-online.target
ConditionPathExists=!/var/lib/niriblue/flatpak-setup.done

[Service]
Type=oneshot
ExecStart=/usr/libexec/niriblue-flatpak-setup
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
